# etc/local.d/local.start
#!/bin/sh
mkdir -p $XDG_RUNTIME_DIR
chmod 700 $XDG_RUNTIME_DIR
setsid /usr/bin/sway &
